# deploy/pogran-bot.service
[Unit]
Description=Pogran Bot - Telegram бот для отслеживания очереди на границе
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=pogran
Group=pogran
WorkingDirectory=/opt/pogran-bot
Environment="PATH=/opt/pogran-bot/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/opt/pogran-bot/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pogran-bot

# Ограничения ресурсов
LimitNOFILE=65536
MemoryLimit=512M

[Install]
WantedBy=multi-user.target

